# ts/conf/ts_sniffer_outbound.conf
# /etc/init/cloudsimd.conf

description "OSRF cloud simulation platform"
author  "Carlos Agüero <caguero@osrfoundation.org>"

start on runlevel [234]
stop on runlevel [0156]

exec ts_sniffer.py 11.8.0.2 > /tmp/ts_sniffer.log 2>&1

respawn
