# vrc_bandwidth/conf/vrc_bandwidth.conf
# /etc/init/DRC_bandwidth.conf

description "OSRF cloud simulation platform"
author  "Carlos Agüero<caguero@osrfoundation.org>"

start on runlevel [234]
stop on runlevel [0156]

exec vrc_bandwidth.py > /tmp/vrc_bandwidth.log 2>&1

respawn
