# vrc_shaping/conf/vrc_controller_outbound.conf
# /etc/init/cloudsimd.conf

description "OSRF cloud simulation platform"
author  "Carlos Agüero <caguero@osrfoundation.org>"

start on runlevel [234]
stop on runlevel [0156]

exec vrc_controller.py -d tun0 -f 0.25 -cl vrc_current_outbound_latency > /tmp/vrc_controller_outbound.log 2>&1

respawn
